# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
user=root

[program:streamlit]
command=streamlit run app.py --server.port=8502
directory=/app
environment=PATH="/venv/bin:%(ENV_PATH)s"
autostart=true
autorestart=true
stderr_logfile=/var/log/streamlit/error.log
stdout_logfile=/var/log/streamlit/stdout.log
priority=200

[program:api]
command=python api.py
directory=/app
environment=PATH="/venv/bin:%(ENV_PATH)s"
autostart=true
autorestart=true
stderr_logfile=/var/log/api/error.log
stdout_logfile=/var/log/api/stdout.log
priority=300

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/stdout.log
priority=400

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

